# deploy/systemd/laundrypro.service
[Unit]
Description=LaundryPro Django app (Gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/laundrypro
EnvironmentFile=/srv/laundrypro/.env
ExecStart=/srv/laundrypro/.venv/bin/gunicorn \
  --workers 3 \
  --timeout 60 \
  --bind 127.0.0.1:8000 \
  config.wsgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
